# services/syslogd/run
#!/bin/sh

: ${SYSLOG_BUFFER_SIZE:=8192}

exec syslogd -n -C$SYSLOG_BUFFER_SIZE

